# systemd/blinky.service
[Unit]
Description=blinky pacman monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/blinky

[Install]
WantedBy=multi-user.target
